# cookbooks/graphite/templates/graphite-api.socket
[Unit]
Description=graphite-api socket

[Socket]
ListenStream=/run/graphite-api.sock
ListenStream=0.0.0.0:8888

[Install]
WantedBy=sockets.target